# packaging/debian/postinst
#!/bin/sh
# foolsboard postinst: create the service user, virtualenv, config + secrets,
# auto-provision a local PostgreSQL database (the default), run migrations, then
# enable and start the systemd service.
set -e

APP_DIR=/opt/foolsboard
DATA_DIR=/var/lib/foolsboard
CONF_DIR=/etc/foolsboard
ENV_FILE="$CONF_DIR/foolsboard.env"
VENV="$APP_DIR/venv"
SVC_USER=foolsboard
DB_NAME=foolsboard
DB_USER=foolsboard

case "$1" in
  configure)
    # 1. Service user (system account, no login shell).
    if ! getent passwd "$SVC_USER" >/dev/null; then
      adduser --system --group --home "$APP_DIR" --no-create-home --quiet "$SVC_USER"
    fi

    # 2. Data + config directories.
    mkdir -p "$DATA_DIR/storage" "$CONF_DIR"
    chown -R "$SVC_USER":"$SVC_USER" "$DATA_DIR"

    # 3. Config + secrets, and a local Postgres database -- generated once on a
    #    fresh install and preserved across upgrades. (An existing config is left
    #    untouched: change DATABASE_URL there by hand to repoint the database.)
    if [ ! -f "$ENV_FILE" ]; then
      systemctl start postgresql >/dev/null 2>&1 || true
      DB_PASS="$(python3 -c 'import secrets; print(secrets.token_urlsafe(24))')"
      if runuser -u postgres -- psql -tAc "SELECT 1 FROM pg_roles WHERE rolname='$DB_USER'" 2>/dev/null | grep -q 1; then
        runuser -u postgres -- psql -c "ALTER ROLE \"$DB_USER\" LOGIN PASSWORD '$DB_PASS'"
      else
        runuser -u postgres -- psql -c "CREATE ROLE \"$DB_USER\" LOGIN PASSWORD '$DB_PASS'"
      fi
      if ! runuser -u postgres -- psql -tAc "SELECT 1 FROM pg_database WHERE datname='$DB_NAME'" 2>/dev/null | grep -q 1; then
        runuser -u postgres -- psql -c "CREATE DATABASE \"$DB_NAME\" OWNER \"$DB_USER\""
      fi

      JWT="$(python3 -c 'import secrets; print(secrets.token_urlsafe(48))')"
      cat > "$ENV_FILE" <<EOF
# foolsboard configuration -- generated on install.

# Where the app listens. By default it serves the SPA + API + WebSocket + media
# itself (no other software needed). To run behind a reverse proxy, point the
# proxy at http://<this-host>:9534 (see /usr/share/doc/foolsboard/).
FOOLSBOARD_HOST=0.0.0.0
FOOLSBOARD_PORT=9534

# Built frontend, served by the app.
STATIC_DIR=/opt/foolsboard/frontend/dist

# Database: a local PostgreSQL auto-provisioned on this host. Point this at an
# external server (or sqlite:////var/lib/foolsboard/foolsboard.db) to change it.
DATABASE_URL=postgresql+psycopg://$DB_USER:$DB_PASS@localhost:5432/$DB_NAME

STORAGE_LOCAL_DIR=/var/lib/foolsboard/storage
JWT_SECRET=$JWT
EOF
      chown root:"$SVC_USER" "$ENV_FILE"
      chmod 640 "$ENV_FILE"
    fi

    # 4. Python virtualenv + dependencies (incl. the psycopg Postgres driver).
    if [ ! -x "$VENV/bin/python" ]; then
      python3 -m venv "$VENV"
    fi
    "$VENV/bin/pip" install --quiet --upgrade pip
    "$VENV/bin/pip" install --quiet -r "$APP_DIR/backend/requirements.txt"
    "$VENV/bin/pip" install --quiet -r "$APP_DIR/backend/requirements-postgres.txt"
    chown -R "$SVC_USER":"$SVC_USER" "$APP_DIR"

    # 5. Database migrations (as the service user).
    runuser -u "$SVC_USER" -- sh -c "set -a; . '$ENV_FILE'; cd '$APP_DIR/backend'; exec '$VENV/bin/alembic' upgrade head"

    # 6. systemd service.
    systemctl daemon-reload
    systemctl enable foolsboard.service >/dev/null 2>&1 || true
    systemctl restart foolsboard.service
    ;;
esac

exit 0
