# packaging/systemd/foolsboard-backup.service
[Unit]
Description=foolsboard backup (restic snapshot of database + media)
After=postgresql.service network.target
Wants=postgresql.service

[Service]
Type=oneshot
# Runs as the foolsboard user so the restic repo has a single owner (the app's
# "Run backup now" runs as this user too). It can read the env (640
# root:foolsboard), read the media dir it owns, reach Postgres via DATABASE_URL,
# and write to the group-writable /var/backups/foolsboard.
User=foolsboard
Group=foolsboard
ExecStart=/opt/foolsboard/backup.sh
